# scripts/mt5.service
[Unit]
Description=MetaTrader 5 (Wine/Headless) for OpenClaw
After=network.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=root
Environment=DISPLAY=:99
Environment=WINEPREFIX=/root/.mt5
Environment=HOME=/root

# Путь к MT5 terminal64.exe в Wine prefix
# Обновите если MT5 установлен в другую папку
ExecStart=/usr/bin/wine "/root/.mt5/drive_c/Program Files/MetaTrader 5/terminal64.exe" /portable

# Ожидаем готовности сети и дисплея
ExecStartPre=/bin/sleep 3

Restart=on-failure
RestartSec=30
TimeoutStopSec=30

# Логирование
StandardOutput=append:/var/log/mt5.log
StandardError=append:/var/log/mt5.log

[Install]
WantedBy=multi-user.target
